# deploy/systemd/pm-robot-paper-settle.timer
[Unit]
Description=Settle pm-robot paper portfolio after each completed run

[Timer]
OnBootSec=5min
OnUnitInactiveSec=10min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
