# deploy/systemd/pm-robot-wallet-history-planner.service
[Unit]
Description=pm-robot bounded wallet history planner
After=network-online.target pm-robot-wallet-level-select.service

[Service]
Type=oneshot
WorkingDirectory=/opt/pm-robot
EnvironmentFile=/opt/pm-robot/.env
TimeoutStartSec=120
Environment=PM_ROBOT_LOCK_WAIT=60
Environment=PM_ROBOT_TASK_NAME=wallet-history-planner
ExecStart=/opt/pm-robot/deploy/scripts/run_locked.sh /opt/pm-robot/.venv/bin/python -m pm_robot.cli --env /opt/pm-robot/.env wallet-history-plan --limit 12 --max-active-jobs 36 --light-refresh-seconds 2592000 --deep-refresh-seconds 604800 --shard-count 3
