# deploy/systemd/pm-robot-trade-role.timer
[Unit]
Description=Run pm-robot official maker/taker evidence ingestion every 30 minutes

[Timer]
OnBootSec=12min
OnUnitInactiveSec=30min
AccuracySec=2min
Persistent=true

[Install]
WantedBy=timers.target
