# deploy/systemd/pm-robot-evidence-worker@.timer
[Unit]
Description=Run pm-robot sharded evidence backfill worker %i

[Timer]
OnBootSec=8min
OnUnitInactiveSec=3min
AccuracySec=20s
Persistent=true
Unit=pm-robot-evidence-worker@%i.service

[Install]
WantedBy=timers.target
