# deploy/systemd/pm-robot-gdrive-backup.service
[Unit]
Description=pm-robot Google Drive streaming SQLite backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/pm-robot
EnvironmentFile=/opt/pm-robot/.env
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
ExecStart=/opt/pm-robot/deploy/scripts/gdrive_backup.sh
