# deploy/systemd/pm-robot-backup.timer
[Unit]
Description=Run pm-robot database backup daily

[Timer]
OnCalendar=*-*-* 04:15:00 UTC
AccuracySec=5min
Persistent=true

[Install]
WantedBy=timers.target
