# deploy/systemd/pm-robot-evidence-planner.timer
[Unit]
Description=Plan pm-robot sharded evidence backfill jobs

[Timer]
OnBootSec=6min
OnUnitInactiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
