# deploy/systemd/pm-robot-wallet-history-planner.timer
[Unit]
Description=Plan bounded pm-robot light/deep history jobs

[Timer]
OnBootSec=4min
OnUnitInactiveSec=3min
AccuracySec=20s
Persistent=true

[Install]
WantedBy=timers.target
