# deploy/systemd/pm-robot-activity.service
[Unit]
Description=pm-robot active wallet activity refresh
After=network-online.target pm-robot-ingest.service

[Service]
Type=oneshot
WorkingDirectory=/opt/pm-robot
EnvironmentFile=/opt/pm-robot/.env
TimeoutStartSec=900
Environment=PM_ROBOT_LOCK_WAIT=30
Environment=PM_ROBOT_TASK_NAME=activity
ExecStart=/opt/pm-robot/deploy/scripts/run_locked.sh /opt/pm-robot/.venv/bin/python -m pm_robot.cli --env /opt/pm-robot/.env ingest-activity --wallet-limit 10 --page-limit 50 --max-events-per-wallet 50 --sleep 0.1
