# deploy/systemd/pm-robot-paper-settle.service
[Unit]
Description=pm-robot paper portfolio settlement
After=network-online.target pm-robot-paper.service pm-robot-gamma.service

[Service]
Type=oneshot
WorkingDirectory=/opt/pm-robot
EnvironmentFile=/opt/pm-robot/.env
TimeoutStartSec=240
Environment=PM_ROBOT_LOCK_WAIT=120
Environment=PM_ROBOT_TASK_NAME=paper-settle
ExecStart=/opt/pm-robot/deploy/scripts/run_locked.sh /opt/pm-robot/.venv/bin/python -m pm_robot.cli --env /opt/pm-robot/.env paper-settle
