# deploy/systemd/pm-robot-materialize-features.service
[Unit]
Description=pm-robot small-batch wallet feature materializer
After=network-online.target pm-robot-evidence-planner.service pm-robot-paper-settle.service

[Service]
Type=oneshot
WorkingDirectory=/opt/pm-robot
EnvironmentFile=/opt/pm-robot/.env
TimeoutStartSec=420
Environment=PM_ROBOT_LOCK_WAIT=30
Environment=PM_ROBOT_TASK_NAME=materialize-features
ExecStart=/opt/pm-robot/deploy/scripts/run_locked.sh /opt/pm-robot/.venv/bin/python -m pm_robot.cli --env /opt/pm-robot/.env materialize-features --limit 40 --min-activity-events 25
